# pkg/delphinusdnsd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/delphinusdnsd"
daemon_flags="-f /var/delphinusdnsd/etc/delphinusdns.conf"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
